# configs/systemd/slo-metrics.service
# Suburban-SOC — compute SOC SLO metrics and index them to soc-slo-metrics (the index
# the SLO dashboard reads). Triggered by slo-metrics.timer every 15 min. This closes the
# "monitoring-of-the-monitor" gap where the dashboard silently froze because the job was
# never scheduled (and a host run failed on the container-only ES_CA path).
#
# Install:
#   sudo cp configs/systemd/slo-metrics.{service,timer} /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now slo-metrics.timer
# Run once now / check:  sudo systemctl start slo-metrics.service ; journalctl -u slo-metrics -n 30
[Unit]
Description=Suburban-SOC SLO metrics -> soc-slo-metrics (SLO dashboard data)
Documentation=https://github.com/voltron-1/Suburban_SOC
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# Runs as the repo owner so it can read scripts/setup/.env and use the user's Docker
# Desktop context for the CA extraction below.
User=tjlam
# audit #167 / NIST AC-6: slo_metrics.py needs READ on .alerts-security.alerts-*,
# soar-actions-*, logstash-security-* plus WRITE to soc-slo-metrics — it now runs as
# the least-privilege `slo_metrics` user (role: slo_metrics_reader), not the `elastic`
# superuser. EnvironmentFile loads SLO_METRICS_PASSWORD from .env; ES_PASS below
# resolves it via systemd's ${VAR} expansion (systemd >= 246).
EnvironmentFile=/home/tjlam/projects/Suburban-SOC/scripts/setup/.env
Environment=ES_URL=https://localhost:9200
Environment=ES_USER=slo_metrics
Environment=ES_PASS=${SLO_METRICS_PASSWORD}
Environment=ES_CA=/run/suburban-soc-slo/ca.crt
# /run/suburban-soc-slo, owned by the service user; auto-removed when the unit stops.
RuntimeDirectory=suburban-soc-slo
# The stack CA lives only in a docker volume (no stable host path under Docker Desktop/WSL),
# so pull it from the ES container each run. Idempotent.
ExecStartPre=/usr/bin/docker cp elasticsearch:/usr/share/elasticsearch/config/certs/ca/ca.crt /run/suburban-soc-slo/ca.crt
ExecStart=/usr/bin/python3 /home/tjlam/projects/Suburban-SOC/scripts/setup/ai_agent/slo_metrics.py
# slo_metrics.py exits 2 when an SLO is breaching — that's a successful run reporting a
# breach, not a job failure. Treat 0 and 2 as success so the timer isn't marked failed.
SuccessExitStatus=0 2

# --- Sandboxing (audit #167 / NIST CM-7) ---------------------------------------
# This unit only makes outbound HTTPS calls (requests) and reads the repo/.env —
# it needs no special Linux capabilities, kernel access, or namespace privileges.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true
ProtectClock=true
ProtectKernelLogs=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectControlGroups=true
ProtectHostname=true
RestrictNamespaces=true
RestrictSUIDSGID=true
RestrictRealtime=true
LockPersonality=true
RemoveIPC=true
UMask=0077
CapabilityBoundingSet=
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
