# docs/detections/attack-coverage.json
{
  "name": "Suburban-SOC Detection Coverage",
  "versions": {
    "attack": "14",
    "navigator": "4.9.1",
    "layer": "4.5"
  },
  "domain": "enterprise-attack",
  "description": "Suburban-SOC ATT&CK coverage (WS1.5). Green = a detection exists (Sigma->Elastic Detection Engine, or a Zeek network detection). Generated by scripts/setup/build_attack_coverage.py.",
  "sorting": 3,
  "hideDisabled": false,
  "techniques": [
    {
      "techniqueID": "T1105",
      "tactic": "command-and-control",
      "score": 100,
      "color": "#2ca02c",
      "comment": "Malicious File Download via Bitsadmin \u2014 rules/sigma/proc_creation_win_bitsadmin_download.yml (test: Detections CI (sigma conversion) + Elastic Detection Engine)",
      "enabled": true
    },
    {
      "techniqueID": "T1070.001",
      "tactic": "defense-evasion",
      "score": 100,
      "color": "#2ca02c",
      "comment": "Clearing Windows Event Logs via Wevtutil \u2014 rules/sigma/proc_creation_win_clear_event_logs.yml (test: Detections CI (sigma conversion) + Elastic Detection Engine)",
      "enabled": true
    },
    {
      "techniqueID": "T1136.001",
      "tactic": "persistence",
      "score": 100,
      "color": "#2ca02c",
      "comment": "Local User Account Creation via Net.exe \u2014 rules/sigma/proc_creation_win_local_acct_create.yml (test: Detections CI (sigma conversion) + Elastic Detection Engine)",
      "enabled": true
    },
    {
      "techniqueID": "T1003.001",
      "tactic": "credential-access",
      "score": 100,
      "color": "#2ca02c",
      "comment": "LSASS Memory Dump via Comsvcs.dll \u2014 rules/sigma/proc_creation_win_lsass_dump.yml (test: Detections CI (sigma conversion) + Elastic Detection Engine)",
      "enabled": true
    },
    {
      "techniqueID": "T1218.005",
      "tactic": "defense-evasion",
      "score": 100,
      "color": "#2ca02c",
      "comment": "Mshta Remote or Script Payload Execution \u2014 rules/sigma/proc_creation_win_mshta_remote.yml (test: Detections CI (sigma conversion) + Elastic Detection Engine)",
      "enabled": true
    },
    {
      "techniqueID": "T1059.001",
      "tactic": "execution",
      "score": 100,
      "color": "#2ca02c",
      "comment": "Suspicious PowerShell Encoded Command Execution \u2014 rules/sigma/proc_creation_win_powershell_encoded.yml (test: Detections CI (sigma conversion) + Elastic Detection Engine)",
      "enabled": true
    },
    {
      "techniqueID": "T1574",
      "tactic": "privilege-escalation",
      "score": 100,
      "color": "#2ca02c",
      "comment": "RDP Session Hijacking via Tscon \u2014 rules/sigma/proc_creation_win_rdp_hijack_tscon.yml (test: Detections CI (sigma conversion) + Elastic Detection Engine)",
      "enabled": true
    },
    {
      "techniqueID": "T1218.010",
      "tactic": "defense-evasion",
      "score": 100,
      "color": "#2ca02c",
      "comment": "Regsvr32 Execution from Remote Server \u2014 rules/sigma/proc_creation_win_regsvr32_remote_sct.yml (test: Detections CI (sigma conversion) + Elastic Detection Engine)",
      "enabled": true
    },
    {
      "techniqueID": "T1053.005",
      "tactic": "execution",
      "score": 100,
      "color": "#2ca02c",
      "comment": "Scheduled Task Creation via Schtasks \u2014 rules/sigma/proc_creation_win_scheduled_task.yml (test: Detections CI (sigma conversion) + Elastic Detection Engine)",
      "enabled": true
    },
    {
      "techniqueID": "T1033",
      "tactic": "discovery",
      "score": 100,
      "color": "#2ca02c",
      "comment": "Suspicious System Owner/User Discovery \u2014 rules/sigma/proc_creation_win_user_discovery.yml (test: Detections CI (sigma conversion) + Elastic Detection Engine)",
      "enabled": true
    },
    {
      "techniqueID": "T1047",
      "tactic": "execution",
      "score": 100,
      "color": "#2ca02c",
      "comment": "WMI Process Call Create \u2014 rules/sigma/proc_creation_win_wmi_process_create.yml (test: Detections CI (sigma conversion) + Elastic Detection Engine)",
      "enabled": true
    },
    {
      "techniqueID": "T1046",
      "tactic": "discovery",
      "score": 100,
      "color": "#2ca02c",
      "comment": "Network Service Discovery \u2014 configs/logstash.conf (Category 5 framework enrichment) (test: tests/pipeline/test_framework_enrichment.py)",
      "enabled": true
    },
    {
      "techniqueID": "T1110",
      "tactic": "credential-access",
      "score": 100,
      "color": "#2ca02c",
      "comment": "Brute Force \u2014 configs/logstash.conf (Category 5 framework enrichment) (test: tests/pipeline/test_framework_enrichment.py)",
      "enabled": true
    }
  ],
  "gradient": {
    "colors": [
      "#ffffff",
      "#2ca02c"
    ],
    "minValue": 0,
    "maxValue": 100
  },
  "legendItems": [
    {
      "label": "Detection deployed",
      "color": "#2ca02c"
    }
  ],
  "metadata": [
    {
      "name": "detections",
      "value": "13"
    },
    {
      "name": "source",
      "value": "rules/sigma/*.yml + configs/logstash.conf"
    }
  ]
}
